# 81-usbbackup.rules
KERNEL=="sd?",SUBSYSTEM=="block",RUN+="/home/osar/bin/usbbackup.sh %s{size}"
